# systemd/enftun@.service.in
[Unit]
Description=ENF tunnel client for %I
Documentation=https://github.com/xaptum/enftun
After=tun@%i.service
BindsTo=tun@%i.service
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/etc/enftun
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/enftun -c %i.conf
ProtectSystem=true
ProtectHome=true
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
